# service/cecdaemon.service
[Unit]
Description=CEC Daemon for HDMI-CEC control
After=multi-user.target

[Service]
Type=simple
User={USER}
Group={GROUP}
WorkingDirectory={WORK_DIR}
ExecStart={EXECUTABLE_PATH}
Environment=CECDAEMON_CONFIG_DIR={CONFIG_DIR}
ConfigurationDirectory=cecdaemon
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={WORK_DIR}
ReadWritePaths={CONFIG_DIR}

[Install]
WantedBy=multi-user.target